# systemd/bequest-heartbeat.service
[Unit]
Description=Bequest heartbeat check (dead man's switch)

[Service]
Type=oneshot
ExecStart=/home/will/.local/bin/bequest heartbeat check
