# systemd/claude-code-nightly-cleanup.service
[Unit]
Description=Kill stale Claude Code processes

[Service]
Type=oneshot
ExecStart=/home/will/.local/bin/claude-code-nightly-cleanup
