# systemd/claude-worktree-cleaner.timer
[Unit]
Description=Daily prune of stale Claude-Code agent git worktrees

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true

[Install]
WantedBy=timers.target
