# systemd/cohs-token-refresh.timer
[Unit]
Description=Refresh COHS M365 Outlook-scope OAuth tokens every 30 min

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
