# systemd/cohs-graph-token-refresh.timer
[Unit]
Description=Refresh COHS M365 Graph-scope OAuth tokens every 30 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
