# systemd/cohs-graph-token-refresh.service
[Unit]
Description=Refresh COHS M365 Graph-scope OAuth tokens (/me/sendMail)
After=network-online.target gnome-keyring-daemon.service
Wants=network-online.target gnome-keyring-daemon.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/practiceforge email m365-refresh
