# systemd/cross-machine-sync-check.timer
[Unit]
Description=Cross-machine sync check (Thu 18:00 + Mon 20:00)
Requires=cross-machine-sync-check.service

[Timer]
OnCalendar=Thu *-*-* 18:00:00
OnCalendar=Mon *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
