# systemd/dev-catchup.service
[Unit]
Description=Catch up on missing dev:: DayPage entries
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dev-catchup --apply
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h

[Install]
WantedBy=default.target
