# systemd/firefox-nightly-restart.service
[Unit]
Description=Gracefully restart Firefox to prevent memory leaks
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -TERM firefox
SuccessExitStatus=1
# Firefox saves session on SIGTERM, so tabs will restore on next launch
# pkill exits 1 when no process matched — that's fine, not a failure

[Install]
WantedBy=default.target
