# systemd/frecency-daemon.timer
[Unit]
Description=Run frecency daemon daily at 4am

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
