# systemd/gmail-token-refresh.timer
[Unit]
Description=Refresh Gmail OAuth access token every 30 min

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
