# systemd/gmail-push-tags.timer
[Unit]
Description=Push Gmail tag changes every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
