# systemd/helix-undo-cleanup.timer
[Unit]
Description=Weekly cleanup of old Helix undo files

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
