# systemd/mail-sync.timer
[Unit]
Description=Run mail-sync every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
