# systemd/mailforge.service
[Unit]
Description=mailforge — Axum daemon for browser-handoff mail rendering
# pizauth.service must be up first: mailforge shells out to `pizauth show`
# for OAuth tokens on every send. network-online.target is transitive via
# pizauth but stated explicitly for clarity.
After=pizauth.service network-online.target
Wants=pizauth.service network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mailforge daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
StandardOutput=append:%h/.local/share/mailforge.log
StandardError=append:%h/.local/share/mailforge.log

[Install]
WantedBy=default.target
