# systemd/mbsync-cohs.service
[Unit]
Description=Pull COHS M365 mail via mbsync (XOAUTH2) and re-index with notmuch
After=network-online.target pizauth.service
Wants=network-online.target pizauth.service

[Service]
Type=oneshot
# mbsync's PassCmd runs `pizauth show cohs`. The pizauth daemon holds the
# OAuth tokens in memory and refreshes transparently — no keyring / D-Bus
# involvement on the sync path. Absolute paths because systemd's default
# PATH is minimal (no ~/.local/bin).
ExecStart=/bin/sh -c '/usr/bin/mbsync cohs && /usr/bin/notmuch new'
