# systemd/meliview.service
[Unit]
Description=meliview — Axum daemon for browser-handoff mail rendering
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/meliview daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
StandardOutput=append:%h/.local/share/meliview.log
StandardError=append:%h/.local/share/meliview.log

[Install]
WantedBy=default.target
