# systemd/package-list-backup.timer
[Unit]
Description=Daily pacman package list backup

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
