# systemd/pizauth-dump.timer
[Unit]
Description=Dump pizauth state every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
