# systemd/state-capture.service
[Unit]
Description=Capture system state (packages, services, crates, groups) to dotfiles

[Service]
Type=oneshot
Environment=PATH=/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:%h/.cargo/bin:%h/.local/bin
ExecStart=%h/.local/bin/state-capture capture

[Install]
WantedBy=default.target
