# systemd/system-health-check.service
[Unit]
Description=Daily system health check

[Service]
Type=oneshot
ExecStart=/home/will/.local/bin/system-health-check
