# systemd/zotero-pdf-watcher.service
[Unit]
Description=Zotero PDF Watcher - Native Nushell Re-nu Version
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/zotero-pdf-watcher-renu
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment="PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=%h"

[Install]
WantedBy=default.target