# systemd/zotero-pdf-watcher.service
[Unit]
Description=Zotero PDF Watcher - Native Nushell Re-nu Version
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/nu /home/%u/.local/bin/zotero-pdf-watcher-renu
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment for Nushell
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/%u"

[Install]
WantedBy=default.target