# macos/LaunchAgents/com.williamnapier.clinic-attendance-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.clinic-attendance-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.cargo/bin/clinic-attendance-report</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/clinic-attendance-report.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clinic-attendance-report.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
    </dict>
</dict>
</plist>
