# macos/LaunchAgents/com.williamnapier.mailforge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.mailforge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/mailforge</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/mailforge.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/mailforge.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/williamnapier/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
        <key>RUST_LOG</key>
        <string>info</string>
        <key>MAILFORGE_STATIC_DIR</key>
        <string>/Users/williamnapier/Code/mailforge/static</string>
    </dict>

    <!-- Interactive, NOT Background (2026-07-06). This is a mail UI the
         user actively waits on. Background QoS pinned the daemon AND every
         notmuch subprocess it spawns to E-cores + throttled I/O — a
         measured ~4.3x penalty per notmuch call (10ms -> 43ms) plus
         multi-second starvation under any foreground CPU load. Interactive
         runs it (and its query children) at normal app QoS. -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
