# macos/launchd/com.williamnapier.mbsync-gmail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.mbsync-gmail</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>/Users/williamnapier/.local/bin/mbsync gmail &amp;&amp; /opt/homebrew/bin/notmuch new</string>
    </array>
    <!-- Permanent Gmail pull. Replaces the transitional mbsync-gmail-spike
         (4h cadence, 2026-05-01 to 2026-05-16). The spike was sized for the
         initial backlog under Workspace's per-hour IMAP rate limit; with the
         backlog drained (steady-state delta ~50-140 msgs/day across the
         last 10 days of MESSAGEBOARD briefings) the rate limit is no
         longer a concern, and 15-min cadence matches the cohs agent.

         EnvironmentVariables block is REQUIRED for the same reason as
         mbsync-cohs.plist: `notmuch new` triggers the post-new hook
         (~/Mail/.notmuch/hooks/post-new), which makes bare `notmuch tag`
         calls that inherit launchd's PATH-less env and would otherwise
         fail with "notmuch: command not found". -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/williamnapier/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
        <key>NOTMUCH_CONFIG</key>
        <string>/Users/williamnapier/Mail/.notmuch-config</string>
    </dict>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/mbsync-gmail.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/mbsync-gmail.log</string>
</dict>
</plist>
