# macos/launchd/com.williamnapier.mbsync-cohs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.mbsync-cohs</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>/Users/williamnapier/.local/bin/mbsync cohs &amp;&amp; /opt/homebrew/bin/notmuch new</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/mbsync-cohs.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/mbsync-cohs.log</string>
</dict>
</plist>
