# macos/launchd/com.williamnapier.forge-link-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.forge-link-manager</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/nu</string>
        <string>-c</string>
        <string>/Users/williamnapier/.local/bin/link-service start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/williamnapier</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/williamnapier/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/williamnapier/scripts/forge-link-manager/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/scripts/forge-link-manager/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>