# macos/launchd/com.williamnapier.m365-imap-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.m365-imap-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/practiceforge</string>
        <string>email</string>
        <string>m365-imap-refresh</string>
    </array>
    <!-- Keep-warm: the PassCmd path refreshes-on-read automatically when
         the token is within 5 min of expiry, so this timer exists purely
         to keep the keychain's access_token fresh between mbsync runs
         (reduces the chance of the very first mbsync pull after idle
         hitting the refresh HTTP call inline). Hourly is sufficient:
         Microsoft access tokens last 60–90 min. -->
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/m365-imap-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/m365-imap-refresh.log</string>
</dict>
</plist>
