# macos/launchd/com.williamnapier.cohs-graph-token-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.cohs-graph-token-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/practiceforge</string>
        <string>email</string>
        <string>m365-refresh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/cohs-graph-token-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/cohs-graph-token-refresh.log</string>
</dict>
</plist>
