# macos/launchd/com.williamnapier.git-auto-commit-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.git-auto-commit-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/dotfiles/scripts/git-auto-commit-watcher-macos</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/williamnapier/.local/share/git-auto-commit-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/.local/share/git-auto-commit-watcher-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/williamnapier/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>
</dict>
</plist>