# macos/launchd/com.williamnapier.gmail-push-tags.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.gmail-push-tags</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/practiceforge</string>
        <string>email</string>
        <string>gmail-push-tags</string>
        <string>--push</string>
    </array>
    <!-- The push-tags half of the Gmail IMAP+notmuch stack: local
         notmuch tag changes get reflected to Gmail as label changes
         every 15 min. State file:
         ~/.local/share/practiceforge/gmail-push-state.json.
         First run bookmarks the current notmuch lastmod and does NO
         API writes — subsequent runs push only deltas. Coexists with
         lieer's mail-sync (which pulls and indexes); lieer's own
         `gmi push` is disabled via PUSH_VIA_GMI=0 in mail-sync. -->
    <!-- EnvironmentVariables is REQUIRED: practiceforge shells out
         to `notmuch` to read tag state, and launchd processes don't
         inherit the interactive PATH. Without this block, every tick
         crashes with "No such file or directory" before any API call.
         Discovered 2026-04-23 after ~20 silent-no-op ticks. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/williamnapier/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
        <key>NOTMUCH_CONFIG</key>
        <string>/Users/williamnapier/Mail/.notmuch-config</string>
    </dict>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/williamnapier/Library/Logs/gmail-push-tags.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/Library/Logs/gmail-push-tags.log</string>
</dict>
</plist>
