# macos/launchd/com.williamnapier.forum-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.williamnapier.forum-worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/forum</string>
        <string>worker</string>
        <string>--poll-seconds</string>
        <string>10</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/williamnapier/Assistants/shared/design-forum</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/williamnapier/.local/share/forum-worker.stdout</string>

    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/.local/share/forum-worker.stderr</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/williamnapier/.local/bin:/Users/williamnapier/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
