# launchd/com.user.mail-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mail-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamnapier/.local/bin/mail-sync</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/williamnapier/.local/share/mail-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/williamnapier/.local/share/mail-sync.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/williamnapier/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/williamnapier</string>
        <key>NOTMUCH_CONFIG</key>
        <string>/Users/williamnapier/Mail/.notmuch-config</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>
</dict>
</plist>
