# linux/systemd/cohs-token-refresh.service
[Unit]
Description=Refresh COHS M365 Outlook-scope OAuth tokens (IMAP/SMTP)
After=network-online.target gnome-keyring-daemon.service
Wants=network-online.target gnome-keyring-daemon.service

[Service]
Type=oneshot
# cohs-oauth is a bash wrapper that invokes the venv Python; need login shell
# PATH so "cohs-oauth" resolves (~/.local/bin) and secret-tool is on PATH.
ExecStart=%h/.local/bin/cohs-oauth refresh
# secret-tool needs DBUS_SESSION_BUS_ADDRESS; for a user unit launched by
# systemd --user this is inherited from the session. If refresh fails with
# "cannot autolaunch D-Bus" it's because the user systemd instance wasn't
# started from a D-Bus session — check with `systemctl --user show-environment`.
